# deploy/hetzner/trustedoss-backup.service
# TrustedOSS Portal — daily local backup (Hetzner demo SaaS).
#
# Wraps scripts/backup.sh: pg_dump (gzip) + workspace tar + manifest.json into
# backups/<UTC-stamp>/, pruning sets older than BACKUP_RETENTION_DAYS (7).
# Local backup is always taken. The ExecStartPost below then pushes the latest
# set offsite via scripts/backup-offsite.sh — a SAFE NO-OP unless
# BACKUP_OFFSITE_REMOTE is set in .env (so this unit works unchanged for
# local-only deployments). See docs/operator-runbook-hetzner.md §11 to enable it.
#
# Paired timer: trustedoss-backup.timer.
[Unit]
Description=TrustedOSS daily local backup (pg_dump + workspace)
Documentation=https://github.com/trustedoss/trusca/blob/main/docs/operator-runbook-hetzner.md
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=trustedoss
WorkingDirectory=/opt/trustedoss/portal
ExecStart=/usr/bin/env bash scripts/backup.sh
# Offsite push (opt-in): no-op unless BACKUP_OFFSITE_REMOTE is set in .env.
ExecStartPost=/usr/bin/env bash scripts/backup-offsite.sh
TimeoutStartSec=1800
